# role/api/files/api_homelab.service
[Unit]
Description=Homelab API service (FastAPI)
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/api_homelab
ExecStart=/opt/api_homelab/api_venv/bin/uvicorn api_homelab:app --host 0.0.0.0 --port 443 --ssl-keyfile /etc/ssl/private/api.homelab.lan.key --ssl-certfile /etc/ssl/certs/api.homelab.lan.crt
Restart=always

[Install]
WantedBy=multi-user.target
