# role/build/files/homelab-blackbox-exporter_arm64/debian/homelab-blackbox-exporter.prerm
#!/bin/sh
set -e

if systemctl is-active --quiet homelab-blackbox-exporter.service; then
    systemctl stop homelab-blackbox-exporter.service
fi

if systemctl is-enabled --quiet homelab-blackbox-exporter.service; then
    systemctl disable homelab-blackbox-exporter.service
fi

if id "blackbox" >/dev/null 2>&1; then
    userdel blackbox
fi

systemctl daemon-reload

exit 0
