# role/proxmox/files/snapshot_vms.timer
[Unit]
Description=Run snapshot_vms.service weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
