# role/main/files/github_pull.timer
[Unit]
Description=Timer to launch github_pull.service every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
