# role/main/files/github_commits.timer
[Unit]
Description=Timer to launch github_commits.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
